# blds/pipewire/pkg.install
function postinst() {
    systemctl --global enable pipewire.socket
}

function postrm() {
    systemctl --global disable pipewire.socket
}
