# blds/pulseaudio/pkg.install
function postinst() {
    systemctl --global enable pulseaudio.socket
    glib-compile-schemas /usr/share/glib-2.0/schemas/
}

function postup() {
    postinst
}

function postrm() {
    systemctl --global disable pulseaudio.socket
}
