# blds/pulseaudio/pkg.install
function postinst() {
	systemctl --global enable pulseaudio.socket
}

function postup() {
	postinst
}
